# XNUCA2019Qualifier/Web/blog_revenge/blog/crontab
*/1 * * * * python /tmp/robot.py >> /var/log/cron.log 2>&1